# templates/consul.upstart.erb
# Consul Agent (Upstart unit)
description "Consul Agent"
start on (local-filesystems and net-device-up IFACE!=lo)
stop on runlevel [06]

env CONSUL=<%= scope.lookupvar('consul::bin_dir') %>/consul
env CONFIG=<%= scope.lookupvar('consul::config_dir') %>


script
    # read settings like GOMAXPROCS from "/etc/default/consul", if available.
    [ -e /etc/default/consul ] && . /etc/default/consul

    export GOMAXPROCS=${GOMAXPROCS:-2}
    exec sudo -u <%= scope.lookupvar('consul::user') %> -g <%= scope.lookupvar('consul::group') %> $CONSUL agent -config-dir $CONFIG <%= scope.lookupvar('consul::extra_options') %>
end script

respawn
respawn limit 10 10
kill timeout 10
